# meta-eca/recipes-web/the-thing-system/steward/steward.service.in
[Unit]
Description=Steward for TheThingSystem
After=steward-init.service

[Service]
Type=simple
EnvironmentFile=-/etc/TheThingSystem/config
ExecStart=@the_thing_system_dir@/steward/start-steward
Restart=on-failure

[Install]
WantedBy=multi-user.target
